# crontab
0 0 * * * /repositories/conagua/get_latest_data.sh